# service/tomlin-api-watcher.path
[Path]
PathModified=/home/pi/api/bin

[Install]
WantedBy=multi-user.target
